# cbw-cac.service
[Unit]
Description=CBW Configuration-as-Code (CAC) Cloud-Init Builder
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
ExecStart=/etc/cbw-cac/cbw_cac_daemon.py
Restart=on-failure
RestartSec=5s
User=root
Group=root
NoNewPrivileges=true
ProtectSystem=full
ProtectHome=true
ReadWritePaths=/etc/cbw-cac /var/lib/cbw-cac /var/log/cbw-cac
PrivateTmp=true

[Install]
WantedBy=multi-user.target
